# examples/ghops.service
# Example systemd service file for ghops
# Save this to /etc/systemd/system/ghops.service or ~/.config/systemd/user/ghops.service
# 
# For system-wide installation (requires sudo):
#   sudo cp ghops.service /etc/systemd/system/
#   sudo systemctl enable ghops
#   sudo systemctl start ghops
#
# For user installation:
#   mkdir -p ~/.config/systemd/user
#   cp ghops.service ~/.config/systemd/user/
#   systemctl --user enable ghops
#   systemctl --user start ghops

[Unit]
Description=GitHub Operations Social Media Service
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=YOUR_USERNAME
Group=YOUR_GROUP
WorkingDirectory=/home/YOUR_USERNAME
Environment=PATH=/home/YOUR_USERNAME/miniforge3/envs/ghops-env/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/home/YOUR_USERNAME/miniforge3/envs/ghops-env/bin/ghops service start
Restart=always
RestartSec=60
StandardOutput=journal
StandardError=journal

# Optional: Set resource limits
MemoryLimit=512M
CPUQuota=50%

[Install]
WantedBy=multi-user.target
# For user service, use: WantedBy=default.target
